# tools/build/docker/s6/s6-rc.d/redis/run
#!/bin/sh
exec s6-setuidgid koyomi /usr/bin/redis-server /home/koyomi/lanraragi/tools/build/docker/redis.conf
