# tools/build/docker/s6/services.d/lanraragi/finish
#!/usr/bin/execlineb -S0
# (-t)erminate running services to allow container to restart
s6-svscanctl -t /var/run/s6/services
