# tools/build/docker/s6/services.d/lanraragi/run
#!/usr/bin/execlineb -P
s6-setuidgid koyomi
cd /home/koyomi/lanraragi/
perl /home/koyomi/lanraragi/script/launcher.pl -f /home/koyomi/lanraragi/script/lanraragi
