# tools/build/docker/s6/services.d/lanraragi/run
#!/bin/sh
cd /home/koyomi/lanraragi/
export HOME=/home/koyomi
s6-setuidgid $(id -u koyomi):$(id -g koyomi) perl /home/koyomi/lanraragi/script/launcher.pl -f /home/koyomi/lanraragi/script/lanraragi
