# tools/build/docker/s6/services.d/redis/run
#!/bin/sh
s6-setuidgid $(id -u koyomi):$(id -g koyomi) /usr/bin/redis-server /home/koyomi/lanraragi/tools/build/docker/redis.conf
