# tools/build/docker/s6/services.d/redis/run
#!/usr/bin/execlineb -P
s6-envuidgid -B koyomi:koyomi
/usr/bin/redis-server
