# tools/build/docker/s6/s6-rc.d/lanraragi/finish
#!/bin/sh
kill -15 `cat /home/koyomi/lanraragi/temp/shinobu.pid-s6`
kill -15 `cat /home/koyomi/lanraragi/temp/minion.pid-s6`
sleep 5
