# tools/build/docker/s6/s6-rc.d/lanraragi/finish
#!/bin/sh

pkill --signal SIGINT -F "/home/koyomi/lanraragi/temp/shinobu.pid-s6"
pkill --signal SIGTERM -F "/home/koyomi/lanraragi/temp/minion.pid-s6"

pidwait -F "/home/koyomi/lanraragi/temp/shinobu.pid-s6"
pidwait -F "/home/koyomi/lanraragi/temp/minion.pid-s6"
