# tools/build/docker/s6/s6-rc.d/lanraragi/run
#!/bin/sh
cd /home/koyomi/lanraragi/
export HOME=/home/koyomi
exec s6-setuidgid koyomi perl /home/koyomi/lanraragi/script/launcher.pl -f /home/koyomi/lanraragi/script/lanraragi
